# beacon-alpha-1/beacon.conf
[program:beacon]
command=/home/ubuntu/beacon
priority=100
autostart=true
autorestart=true
