# stacks/bitcoin-private/bitcoin/txgen.conf
[program:txgen]
user=ubuntu
directory=/home/ubuntu
environment = HOME='/home/ubuntu/'
command=python3 /home/ubuntu/txgen.py
priority=100
autostart=true
autorestart=true
