# relay-alpha-1/node_exporter.conf
[program:node_exporter]
command=/home/ubuntu/node_exporter
priority=100
autostart=true
autorestart=true
