# relay-alpha-1/relay.conf
[program:relay]
command=/home/ubuntu/relay
priority=100
autostart=true
autorestart=true
