# orchestration/katara/relay/relay.conf
[program:relay]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/relay "192.168.14.23:8002"
priority=100
autostart=true
autorestart=true
