# orchestration/bitcoin-private/bitcoin/blockgen.conf
[program:blockgen]
user=ubuntu
directory=/home/ubuntu
environment = HOME='/home/ubuntu/'
command=python3 /home/ubuntu/blockgen.py
priority=100
autostart=true
autorestart=true
