# pulumi-trial-2019-10-12-parity-private/ansible/txgen.conf
[program:txgen]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/txgen
priority=100
autostart=false
autorestart=false
