# ethereum_node/geth.conf
[program:geth]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/geth --rinkeby --cache=2048 --datadir=/home/ubuntu/eth_data --nousb
priority=100
autostart=true
autorestart=true
