# ethereum_node/node_exporter.conf
[program:node_exporter]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/node_exporter
priority=100
autostart=true
autorestart=true
