# ethereum_node/onramp.conf
[program:onramp]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/onramp
priority=100
autostart=true
autorestart=true
