# deployment/beacon/files/ibeacon.conf
[program:ibeacon]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/beacon --discovery_addr "0.0.0.0:9002" --heartbeat_addr "0.0.0.0:9003"
priority=100
autostart=true
autorestart=true
