# orchestration/katara/msggen/bridge.conf.j2
[program:bridge]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/bridge -b"{{beacon_ip}}:8002"
priority=100
autostart=true
autorestart=true
