# orchestration/katara/msggen/chain.conf
[program:chain]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/chain --datadir /home/ubuntu/eth_data --nousb --rinkeby
priority=100
autostart=true
autorestart=true
