# stacks/katara/msggen/bot.conf
[program:bot]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/bot
priority=100
autostart=true
autorestart=true
