# orchestration/katara/master/master.conf.j2
[program:master]
user=ubuntu
directory=/home/ubuntu
command=/home/ubuntu/master "{{beacon_ip}}:8002"
priority=100
autostart=true
autorestart=true
